# deploy/systemd/yule.target
[Unit]
Description=Yule engineering runtime — all queue workers + supervisor
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target
