# apps/forgekit-console/examples/runtime/com.forgekit.runtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd — forgekit always-on (macOS). 설치: cp ~/Library/LaunchAgents/ →
     launchctl load ~/Library/LaunchAgents/com.forgekit.runtime.plist
     ⚠ macOS sleep(lid close) 시 프로세스 SUSPEND — README 참조. 상시 운영은 Linux/systemd 권장. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.forgekit.runtime</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>forgekit</string>
    <string>runtime</string>
    <string>serve</string>
    <string>--interval</string>
    <string>300</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict><key>FORGEKIT_NOTIFY</key><string>1</string></dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/forgekit-runtime.out.log</string>
  <key>StandardErrorPath</key><string>/tmp/forgekit-runtime.err.log</string>
</dict>
</plist>
